# dist/trusty/etc/init/flywheel.conf
# flywheel - AWS cost control proxy
#
# flywheel is a http proxy which stops ASG/EC2 instances after an idle period, and
# restarts on request

description "AWS cost control proxy"

start on runlevel [345]
stop on runlevel [0126]

respawn

exec /usr/bin/flywheel --config /etc/flywheel.json
